# src/Haus.Web.Host/haus-web.service
[Unit]
Description=Haus web application
After=network.target

[Service]
Type=notify
WorkingDirectory=/home/bryce/haus/web
ExecStart=/home/bryce/haus/web/Haus.Web.Host "--urls http://*:5000;https://*:5001"
SyslogIdentifier=Haus Web
User=bryce
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target